# s6-rc.d/danted/run
#!/usr/bin/execlineb
with-contenv
s6-setuidgid root

foreground {
    echo "starting Danted socks server"
}

exec /usr/sbin/danted -f /etc/danted.conf