# s6-rc.d/tun0-wait/up
/etc/s6-overlay/s6-rc.d/tun0-wait/wait.sh